# cec2/cec-sleep-daemon.service
[Unit]
Description=CEC Sleep Daemon - Detect sleep/wake and control TV
Documentation=file:///home/server/Projects/scripts/cec2/README.md
After=dbus.service
RefuseManualStart=no
RefuseManualStop=no

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/server/Projects/scripts/cec2/cec-sleep-daemon.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

# Keep service running even during system sleep
KillMode=none
KillSignal=SIGTERM
TimeoutStopSec=0

[Install]
WantedBy=multi-user.target
