# cec2/cec_daemon.service
[Unit]
Description=CEC Daemon - Turn TV/Monitor on/off with KDE Plasma
Documentation=file:///path/to/cec2/README.md
After=dbus.service

[Service]
Type=simple
ExecStart=/path/to/cec2/cec_daemon.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=%h/.cache
ReadWritePaths=%h/.config/systemd/user

[Install]
WantedBy=default.target
